# webhook-app/webhook.container.example
[Unit]
Description=Quadlet Deploy Webhook
After=network-online.target
Wants=network-online.target

[Container]
Image=ghcr.io/org/quadlet-webhook:latest
ContainerName=quadlet-webhook
# Nginx reverse proxy kullanıldığında container sadece localhost'a açılır.
PublishPort=127.0.0.1:18080:8080
Environment=SALT_SECRET=replace_me
Environment=VERSION_FILE=/data/global_version
Environment=TZ=UTC
Environment=TOKEN_TOLERANCE_MINUTES=5
Environment=MAX_HEADER_VALUE_LEN=128
Volume=/opt/quadlet-rollout:/data:Z
NoNewPrivileges=true
ReadOnly=true

[Service]
Restart=always
RestartSec=2
ReadWritePaths=/data

[Install]
WantedBy=multi-user.target
