# examples/server-quadlets/users/appuser2/appsvc.container
[Unit]
Description=App Service 2
After=network-online.target

[Container]
Image=ghcr.io/org/appsvc:latest
ContainerName=appsvc
AutoUpdate=registry
EnvironmentFile=%h/.config/quadlet-agent/app.env

[Service]
Restart=always

[Install]
WantedBy=default.target
