# nginx-rollout/systemd/nginx-rollout.service
[Unit]
Description=Nginx + Certbot Rollout Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=NGINX_ROLLOUT_CONFIG=__NGINX_ROLLOUT_ENV_PATH__
ExecStart=/usr/local/bin/nginx-rollout.sh
